# packages/9router-bin/9router.service
[Unit]
Description=9Router AI gateway
Documentation=https://github.com/decolua/9router
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=HOSTNAME=127.0.0.1
EnvironmentFile=-%h/.9router/.env
ExecStart=/usr/bin/9router --no-browser --host 127.0.0.1
Restart=on-failure
RestartSec=5
WorkingDirectory=%h

# Security hardening
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.9router /tmp
NoNewPrivileges=true

[Install]
WantedBy=default.target
